# opdockerfile/mysql/centos/mha_node/preparation/conf.d/mysql.ini
[program:mysql]
directory=/preparation/
command=python /preparation/mysql.py 
user=root
autorestart=true
